# base/etc/cron.hourly/9needs-restarting.cron
#!/bin/bash
needs-restarting -r > /dev/null || shutdown -r
